# recipes-core/entropy-tss/init-api_key_tdx
### BEGIN INIT INFO
# Provides:             api_key_tdx
# Required-Start:       $remote_fs $syslog $networking
# Required-Stop:        $remote_fs $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         1
# Short-Description:    Start and stop the api_key_tdx daemon
### END INIT INFO
#

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/usr/bin/api_key_tdx
NAME=api_key_tdx_vc
DESC="API Key Service"
PIDFILE=/var/run/api_key_tdx_vc.pid


start() {
        echo -n "Starting $DESC: "
        mount -o remount,size=90% /var/volatile
        # This is a mock placeholder for the persistent fs until https://github.com/entropyxyz/meta-entropy-tss/pull/5 is merged
        mkdir -p /persist
        start-stop-daemon -S -p $PIDFILE -N -10 -b -a /bin/sh -- -c "exec ${DAEMON} \
                --box-url 0.0.0.0:3001 \
                2>&1 | tee /tmp/api_key_tdx-logs"
        echo "$NAME."
}

stop() {
        echo -n "Stopping $DESC: "
        start-stop-daemon -K -x "$DAEMON" -p $PIDFILE
        echo "$NAME."
}
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart|reload)
        stop
        start
        ;;
  *)
        N=/etc/init.d/$NAME
        echo "Usage: $N {start|stop|restart|reload}" >&2
        exit 1
        ;;
esac
exit 0
