# recipes-core/entropy-tss/init
### BEGIN INIT INFO
# Provides:             entropy-tss
# Required-Start:       $remote_fs $syslog $networking
# Required-Stop:        $remote_fs $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         1
# Short-Description:    Start and stop the entropy-tss daemon
### END INIT INFO
#

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/usr/bin/entropy-tss
NAME=entropy_tss_vc
DESC="Entropy Threshold Signature Server"
PIDFILE=/var/run/entropy_tss_vc.pid


start() {
        echo -n "Starting $DESC: "
        mount -o remount,size=90% /var/volatile
        mkdir -p /persist
        mount -o remount /dev/nvme0n1p2 /persist
        start-stop-daemon -S -p $PIDFILE -N -10 -b -a /bin/sh -- -c "exec ${DAEMON} \
                2>&1 | tee /tmp/entropy-tss-logs"
        echo "$NAME."
}

stop() {
        echo -n "Stopping $DESC: "
        start-stop-daemon -K -x "$DAEMON" -p $PIDFILE
        echo "$NAME."
}
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart|reload)
        stop
        start
        ;;
  *)
        N=/etc/init.d/$NAME
        echo "Usage: $N {start|stop|restart|reload}" >&2
        exit 1
        ;;
esac
exit 0
